# s6/dell-fans/run
#!/usr/bin/env sh
# Dell PowerEdge adaptive fan controller. The script's own vendor guard
# refuses to start on non-Dell BMCs. Here we additionally short-circuit
# when /dev/ipmi0 isn't present (non-IPMI host, e.g. Unraid on consumer
# hardware) so the rest of the host-agent stack still runs.
if [ ! -c /dev/ipmi0 ]; then
  echo "[dell-fans] /dev/ipmi0 absent — fan controller disabled on this host"
  exec sleep infinity
fi
exec /usr/local/bin/dell-fan-controller.sh
